# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:planqk-proxy]
command=npm start --prefix ./planqk-proxy/
redirect_stderr=true
